# epg.service
# path:   /home/klassiker/.local/share/repos/epg/epg.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2022-11-19T19:26:44+0100

[Unit]
Description=download epg data in user contaxt alarm

[Service]
Type=oneshot
ExecStart=%h/.local/share/repos/epg/epg.sh

[Install]
WantedBy=multi-user.target
