# epg.service
# path:       /home/klassiker/.local/share/repos/epg/epg.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/epg
# date:       2020-06-02T10:53:51+0200

[Unit]
Description=download epg data in user contaxt alarm

[Service]
User=alarm
Type=oneshot
ExecStart=/home/alarm/bin/epg.sh

[Install]
WantedBy=multi-user.target
