# epg.service
# path:   /home/klassiker/.local/share/repos/epg/epg.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2022-11-28T18:49:54+0100

[Unit]
Description=processing and uploading epg data in user context alarm

[Service]
Type=oneshot
ExecStart=%h/.local/share/repos/epg/epg.sh

[Install]
WantedBy=multi-user.target
