# epg.service
# path:   /home/klassiker/.local/share/repos/epg/epg.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2023-01-20T12:27:27+0100

[Unit]
Description=processing and uploading epg data in user context alarm

[Service]
Type=oneshot
ExecStart=/srv/http/download/epg/epg.sh

[Install]
WantedBy=multi-user.target
