# epg.timer
# path:   /home/klassiker/.local/share/repos/epg/epg.timer
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2023-01-25T10:50:59+0100

[Unit]
Description=daily run epg.service

[Timer]
OnCalendar=5:05:00
AccuracySec=1m

[Install]
WantedBy=timers.target
