# epg.timer
# path:   /home/klassiker/.local/share/repos/epg/epg.timer
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2022-11-15T12:42:36+0100

[Unit]
Description=daily run epg.service

[Timer]
OnCalendar=5:00:00
AccuracySec=15m
Unit=epg.service

[Install]
WantedBy=timers.target
