# epg.timer
# path:   /home/klassiker/.local/share/repos/epg/epg.timer
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/epg
# date:   2024-10-06T06:53:31+0200

[Unit]
Description=run epg.service daily

[Timer]
OnCalendar=4:04:00
AccuracySec=1m

[Install]
WantedBy=timers.target
