# epg.service
# path:       /home/klassiker/.local/share/repos/iptv/epg.service
# author:     klassiker [mrdotx]
# github:     https://github.com/mrdotx/iptv
# date:       2020-05-23T00:21:56+0200

[Unit]
Description=download epg data in user contaxt alarm

[Service]
User=alarm
Type=oneshot
ExecStart=/home/alarm/bin/epg.sh

[Install]
WantedBy=multi-user.target
